# spiritmate-myob-sync/rootfs/etc/services.d/spiritmate-myob-sync/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# SpiritMate MYOB Sync - service runner
# ==============================================================================

# Read configuration from add-on options
IMAP_HOST=$(bashio::config 'imap_host')
IMAP_PORT=$(bashio::config 'imap_port')
IMAP_USER=$(bashio::config 'imap_user')
IMAP_PASS=$(bashio::config 'imap_pass')
IMAP_MAILBOX=$(bashio::config 'imap_mailbox')

FROM_EXACT=$(bashio::config 'from_exact')
SUBJECT_PREFIX=$(bashio::config 'subject_prefix')
LABEL_PROCESSED=$(bashio::config 'label_processed')

FIRESTORE_PROJECT_ID=$(bashio::config 'firestore_project_id')

SCHEDULE_ENABLED=$(bashio::config 'schedule_enabled')
SCHEDULE_CRON=$(bashio::config 'schedule_cron')
LOG_LEVEL=$(bashio::config 'log_level')

# Validate required configuration
if bashio::var.is_empty "${IMAP_HOST}" || \
   bashio::var.is_empty "${IMAP_PORT}" || \
   bashio::var.is_empty "${IMAP_USER}" || \
   bashio::var.is_empty "${IMAP_PASS}" || \
   bashio::var.is_empty "${FROM_EXACT}" || \
   bashio::var.is_empty "${SUBJECT_PREFIX}"; then
  bashio::log.error "Missing required IMAP/email settings. Please configure imap_* and email filters."
  exit 1
fi

# Export basic env vars that server needs (but NOT GOOGLE_APPLICATION_CREDENTIALS yet)
export IMAP_HOST IMAP_PORT IMAP_USER IMAP_PASS IMAP_MAILBOX
export FROM_EXACT SUBJECT_PREFIX LABEL_PROCESSED
export FIRESTORE_PROJECT_ID LOG_LEVEL
export SCHEDULE_ENABLED SCHEDULE_CRON

# Start ingress/status/upload server early WITHOUT credentials env var
# This prevents Google Cloud SDK from auto-loading the credentials file
PORT=8099 node /app/dist/server.js &

# Check Firestore credentials (wait for upload if missing)
if [[ ! -f "/share/spiritmate/service-account.json" ]]; then
  bashio::log.warning "Firebase service account not found at /share/spiritmate/service-account.json"
  bashio::log.warning "Open Web UI and upload service-account.json; waiting..."
  while [[ ! -f "/share/spiritmate/service-account.json" ]]; do
    sleep 5
  done
  bashio::log.info "service-account.json detected. Continuing startup..."
fi

# Copy credentials where Node can read them
cp /share/spiritmate/service-account.json /app/service-account.json

# NOW set the credentials path (after server is already running)
GOOGLE_APPLICATION_CREDENTIALS=/app/service-account.json

# Export credentials env var (other vars already exported above)
export GOOGLE_APPLICATION_CREDENTIALS

bashio::log.info "Starting MYOB sync worker"
bashio::log.info "IMAP: ${IMAP_USER}@${IMAP_HOST}:${IMAP_PORT}/${IMAP_MAILBOX:=INBOX}"
bashio::log.info "Schedule enabled: ${SCHEDULE_ENABLED} (${SCHEDULE_CRON})"

# Server already started above (without GOOGLE_APPLICATION_CREDENTIALS set)

if [[ "${SCHEDULE_ENABLED}" == "true" ]]; then
  # Create cron job with full environment
  cat > /tmp/sync-cron.sh <<EOF
#!/bin/bash
export IMAP_HOST="${IMAP_HOST}"
export IMAP_PORT="${IMAP_PORT}"
export IMAP_USER="${IMAP_USER}"
export IMAP_PASS="${IMAP_PASS}"
export IMAP_MAILBOX="${IMAP_MAILBOX}"
export FROM_EXACT="${FROM_EXACT}"
export SUBJECT_PREFIX="${SUBJECT_PREFIX}"
export LABEL_PROCESSED="${LABEL_PROCESSED}"
export FIRESTORE_PROJECT_ID="${FIRESTORE_PROJECT_ID}"
export GOOGLE_APPLICATION_CREDENTIALS="${GOOGLE_APPLICATION_CREDENTIALS}"
export LOG_LEVEL="${LOG_LEVEL}"
cd /app && node dist/index.js
EOF
  chmod +x /tmp/sync-cron.sh
  echo "${SCHEDULE_CRON} /tmp/sync-cron.sh" > /etc/crontabs/root
  crond -f &
  while true; do
    sleep 300
    pgrep crond >/dev/null || crond -f &
  done
else
  bashio::log.info "Schedule disabled. UI is available; you can trigger runs from the panel. Running worker once now and keeping service alive."
  node /app/dist/index.js || bashio::log.warning "Initial run exited with code $?"
  while true; do
    sleep 3600
  done
fi
